# debian-www/php7-fpm/rootfs/etc/services.d/memcached/run
#!/bin/execlineb -P

exec -c
memcached -u daemon -S -l 127.0.0.1
