# alpine-www/hiawatha/service/hiawatha/run
#!/bin/sh

if [ -f ./setup ]
then
  source ./setup
fi

exec hiawatha -d
