# alpine-webapps/owncloud/hiawatha/service/owncloud/run
#!/bin/sh

if [ -f ./setup ]
then
  source ./setup
fi

chown -R root:nobody /usr/share/webapps/owncloud/config
chown -R root:nobody /var/owncloud
