# alpine-mahara/service/apache/run
#!/bin/sh

if [ -f ./setup ]
then
  source ./setup
fi

# display logs
#tail -F /var/log/apache2/*log &

exec httpd -D FOREGROUND
