# alpine-eleasticsearch/service/elasticsearch/run
#!/bin/sh

exec su -c "/opt/elasticsearch/bin/elasticsearch" elasticsearch
