# alpine-webapps/mahara/service/memcached/run
#!/bin/sh

if [ -f ./setup ]
then
  source ./setup
fi

exec /usr/bin/memcached  -u memcached
