# alpine-webapps/etherpad/service/etherpad/run
#!/bin/sh

if [ -f ./setup ]
then
  source ./setup
fi

#exec su -c "/opt/etherpad-lite/bin/run.sh" etherpad 
exec /opt/etherpad-lite/bin/run.sh --root
