# alpine-webapps/etherpad/service/etherpad/run
#!/bin/sh

if [ -f ./setup ]
then
  source ./setup
fi

exec /usr/share/webapps/etherpad/bin/run.sh --root
