# alpine-webapps/etherpad/service/etherpad/run
#!/bin/sh

if [ -f ./setup ]
then
  source ./setup
fi

exec /opt/etherpad-lite/bin/run.sh --root
