# alpine-www/alpine-redis/service/redis/run
#!/bin/sh

if [ -f ./setup ]
then
  source ./setup
fi

exec s6-setuidgid redis redis-server
