# alpine-db/redis/service/redis/run
#!/bin/sh

if [ -f ./setup ]
then
  source ./setup
fi

exec redis-server /etc/redis.conf
