# strut/conf/supervisord-php.conf
[program:php5-fpm]
#command = /usr/sbin/php5-fpm
command = /start-phpfpm.sh
user = root
autostart = true
